# makefiles/cube-base/docker/docker.path
[Unit]
Description=Docker Flanneld

[Path]
PathExists=/run/flannel/subnet.env
Unit=docker.service
