# build.assets/docker/os-rootfs/lib/systemd/system/kubereserved.slice
[Unit]
Description=Kubernetes pods slice
Before=slices.target

[Slice]
