# build.assets/docker/os-rootfs/lib/systemd/system/kubernetes.slice
[Unit]
Description=Kubernetes system components slice
Before=slices.target

[Slice]

