# build.assets/makefiles/master/k8s-master/kube-scheduler.service
[Unit]
Description=Kubernetes Scheduler
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
Requires=kube-apiserver.service
After=kube-apiserver.service

[Service]
ExecStart=/usr/bin/kube-scheduler --master=https://apiserver:6443 --kubeconfig=/var/state/scheduler.kubeconfig
Restart=always
RestartSec=10
User=planet
Group=planet
SuccessExitStatus=2
