# build.assets/makefiles/master/k8s-master/kube-scheduler.service
[Unit]
Description=Kubernetes Scheduler
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
Requires=kube-apiserver.service
After=kube-apiserver.service

[Service]
ExecStart=/usr/bin/kube-scheduler --master=127.0.0.1:8080
Restart=always
RestartSec=10
User=planet
Group=planet
